# salvium-pool.service
# Salvium Pool Systemd Service
# /etc/systemd/system/salvium-pool.service

[Unit]
Description=Salvium Mining Pool
After=network.target
Requires=network.target

[Service]
Type=simple
User=pooluser
Group=pooluser
WorkingDirectory=/home/pooluser/pool
ExecStart=/home/pooluser/pool/salvium-pool -c /home/pooluser/pool/pool.conf
ExecReload=/bin/kill -HUP $MAINPID
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=salvium-pool
RestartPreventExitStatus=10
Restart=on-failure

# Security settings
#NoNewPrivileges=true
#ProtectSystem=strict
#ProtectHome=true
#ReadWritePaths=/home/pooluser/pool
#PrivateTmp=true
#ProtectKernelTunables=true
#ProtectControlGroups=true
#RestrictRealtime=true
#RestrictNamespaces=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=HOME=/home/pooluser
Environment=USER=pooluser
Environment=SALVIUM_ROOT=/home/pooluser/salvium_io
Environment=SALVIUM_BUILD_ROOT=/home/pooluser/salvium_io/build/Linux/main/release
Environment=MONERO_RANDOMX_FULL_MEM=1
Environment=RANDOMX_FULL_MEM=1
Environment=SALVIUM_RANDOMX_FULL_MEM=1

# CPU scheduling
CPUSchedulingPolicy=1
CPUSchedulingPriority=50
IOSchedulingClass=1
IOSchedulingPriority=4

# Memory optimizations
Environment=MALLOC_ARENA_MAX=2
Environment=MALLOC_MMAP_THRESHOLD_=131072

# Network optimizations
ExecStartPre=/bin/sh -c 'echo 65536 > /proc/sys/net/core/somaxconn'
ExecStartPre=/bin/sh -c 'echo 1 > /proc/sys/net/ipv4/tcp_tw_reuse'

[Install]
WantedBy=multi-user.target

